# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=/app/.venv/bin/uvicorn client.backend.main:app --host 127.0.0.1 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/app",MODEL_BACKEND="%(ENV_MODEL_BACKEND)s",MINIMAX_API_KEY="%(ENV_MINIMAX_API_KEY)s",VLLM_ENDPOINT="%(ENV_VLLM_ENDPOINT)s",REPLAY_SCENARIO_ID="%(ENV_REPLAY_SCENARIO_ID)s",PARQUET_DIR="%(ENV_PARQUET_DIR)s",SEMANTIC_LAYER_PATH="%(ENV_SEMANTIC_LAYER_PATH)s",HF_TOKEN="%(ENV_HF_TOKEN)s"

[program:frontend]
command=node /app/client/frontend/server.js
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT="3000",HOSTNAME="127.0.0.1"
